# systemd/huly-bridge.service
[Unit]
Description=Huly.io Bridge Server
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/huly-bridge --config /etc/huly-bridge/bridge.toml
Restart=on-failure
RestartSec=5s
WatchdogSec=30s

# Security hardening
DynamicUser=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
ReadOnlyPaths=/etc/huly-bridge
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=huly-bridge

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
