# systemd/huly-bridge@.service
[Unit]
Description=Huly bridge for workspace %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/huly-bridge --config /etc/huly-bridge/workspaces.d/%i.toml
Restart=on-failure
RestartSec=2s

# Run as a dedicated dynamic user; one per template instance.
DynamicUser=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
ReadOnlyPaths=/etc/huly-bridge
PrivateTmp=yes

# Logging — instance name (%i) appears in the unit name; keep identifier explicit
# for journalctl --identifier consumers.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=huly-bridge@%i

# Resource limits (per workspace).
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
